# supervisord.conf
[supervisord]
nodaemon=true

[program:user-web]
command=paster serve src/delikat/web/user.ini --reload

[program:worker]
command=python -m delikat.worker
